# systemd/netsentinel-capture.service
[Unit]
Description=NetSentinel Packet Capture Service
After=network.target redis.service
Requires=redis.service
Documentation=https://github.com/secuaas/NetSentinel

[Service]
Type=simple
User=root
Group=netsentinel
ExecStart=/opt/netsentinel/bin/netsentinel-capture --config /opt/netsentinel/config/capture.toml
Restart=always
RestartSec=5

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/netsentinel
PrivateTmp=yes

# Required capabilities for packet capture
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

# Resource limits
LimitNOFILE=65535
LimitMEMLOCK=infinity

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

[Install]
WantedBy=multi-user.target
